# custom_uwf.service
[Unit]
Description=uwf modifier

[Service]
ExecStart=/etc/init.d/mystartup.sh

[Install]
WantedBy=multi-user.target
